# docker/etc/services.d/redis/run
#!/command/with-contenv sh
/usr/bin/redis-server