# docker/etc/services.d/celery/run
#!/command/with-contenv sh
cd /app
pipenv run celery -A app.celery worker